# install/templates/template.service
[Unit]
Description=Jira Dashboard (port ${PORT})
After=network.target

[Service]
Type=simple
WorkingDirectory=${ROOT}
ExecStart=${NODE} ${ROOT}/server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=${PROJECT_DIR}/.env

[Install]
WantedBy=default.target
